# com.setugk.journery-agrams-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.setugk.journery-agrams-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/setugk/.journery-agrams/backup_agrams_mac.sh</string>
    </array>

    <!-- Run at 2:30am daily (staggered from Setu's 2:00am job) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Also run once on login (catches missed runs while Mac was asleep) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/setugk/.journery-agrams/backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/setugk/.journery-agrams/backup.log</string>
</dict>
</plist>
